# deploy/hysteria-reload-on-config-change.service
[Unit]
Description=Reload Hysteria when config changes

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart hysteria-server.service
